# src/terok/resources/systemd/terok-clearance-notifier.service
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0
#
# terok-clearance-notifier-version: {{UNIT_VERSION}}
#
# systemd user service for the terok-clearance-notifier — the process
# that bridges clearance-hub events to freedesktop desktop
# notifications.  Separate from the hub unit (terok-dbus.service) so
# headless hosts can run the hub without pulling in a desktop stack.
#
# The {{BIN}} placeholder is resolved at install time by the terok
# setup hook so the unit works regardless of PATH.
#
# Install:
#   terok setup
#   — or manually —
#   sed "s|{{BIN}}|$(which terok-clearance-notifier)|" \
#       resources/systemd/terok-clearance-notifier.service \
#       > ~/.config/systemd/user/terok-clearance-notifier.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now terok-clearance-notifier
#
# Logs:
#   journalctl --user -u terok-clearance-notifier -f

[Unit]
Description=terok clearance desktop notifier
Documentation=https://github.com/terok-ai/terok
# Needs the hub running before we can subscribe; ``Wants=`` (not
# ``Requires=``) so a crashed hub restarts us with a clean retry rather
# than forcibly tearing us down.
Wants=terok-dbus.service
After=terok-dbus.service
# Freedesktop Notifications lives on the session bus.
After=dbus.socket
Requires=dbus.socket

[Service]
Type=simple
ExecStart={{BIN}}
Restart=on-failure
RestartSec=5

# The notifier only speaks varlink upstream (to the hub) and D-Bus
# downstream (to freedesktop.Notifications) — no podman, no userns
# setns, no shield exec.  The hardening cascade that blocks the hub
# doesn't apply here.  Start conservative and layer tightening in a
# follow-up PR once the basic flow is validated in the wild.

[Install]
WantedBy=default.target
