# src/terok/resources/systemd/terok-gate.socket
# SPDX-FileCopyrightText: 2025-2026 Jiri Vyskocil <jiri@vyskocil.com>
#
# SPDX-License-Identifier: Apache-2.0

# terok-gate-version: {{UNIT_VERSION}}

[Unit]
Description=terok git gate socket

[Socket]
ListenStream=127.0.0.1:{{PORT}}
Accept=yes

[Install]
WantedBy=sockets.target
